# sources/init-radio.service
[Unit]
Description=Simple WebRadio service

[Service]
ExecStart=/home/pi/webradio_control.py

[Install]
WantedBy=multi-user.target
